# docker/morscout/consul.service
[Unit]
Description=Consul Service

[Service]
Type=simple
ExecStart=/bin/bash /root/docker/morscout/consul.sh

[Install]
WantedBy=multi-user.target
